# bestseller/w-offline-compute-cloud/src/main/resources/crontab
0 1 * * * /home/hadoop/wanalytics-offline-compute/run.sh all > /home/hadoop/wanalytics-offline-compute/log
